# deploy/jetson/systemd/warehouse-safety.service
[Unit]
# Emergency Guardian — Layer 1 software safety (doc12:80-84, :383 50ms target):
# AMCL distance watch -> Nav2 goal cancel + cmd_vel stop, battery policy, blocked
# timeout, /emergency/event. MUST be up whenever the robots can move; nav2
# BindsTo= this unit so motion never runs without the guardian (and stops if the
# guardian crashes). Final defence line is Layer 0 on the MCU (doc12:75-78).
Description=Warehouse Emergency Guardian (Layer 1 software safety)
Documentation=https://github.com/lll-kkk-ryuya/miniature-warehouse-robotics/blob/main/docs/setup/jetson-deploy.md
After=warehouse-state-cache.service
Wants=warehouse-state-cache.service
PartOf=warehouse.target
# Disable the start-rate limiter so the GUARDIAN itself keeps restarting after
# repeated crashes (pairs with Restart=on-failure below) — the Layer 1 safety unit
# must never stick in 'failed'. NOTE: this auto-recovery is for the guardian ONLY.
# BindsTo= propagates the guardian's stop/inactivation to the motion units (nav2)
# but NOT a start: once the guardian crashes, nav2 stays stopped (fail-safe) even
# after the guardian recovers, and must be restarted manually (systemd.unit(5)
# BindsTo=; recovery procedure in docs/setup/jetson-deploy.md §7).
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=-/etc/warehouse/warehouse.env
User=warehouse
Group=warehouse
RuntimeDirectory=warehouse
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=yes
ExecStart=/opt/warehouse/deploy/jetson/bin/ros-exec.sh ros2 run warehouse_safety emergency_guardian
Restart=on-failure
RestartSec=2

[Install]
WantedBy=warehouse.target
