# deploy/jetson/systemd/warehouse-nav2.service
[Unit]
# Nav2 stack — composes the per-bot Nav2 + twist_mux + VirtualScan via the
# top-level bring-up (warehouse_bringup/launch/bringup.launch.py; itself Includes
# nav2_bringup.launch.py — nav2 composed by #75/PR#93, full Phase 0.5 stack by #156/PR#162).
# bringup.launch.py now DEFAULTS to sim:=true / llm:=true for the Mac M4 capstone
# (bringup.launch.py:147-148,153-154). prod is real hardware (no Gazebo) and runs the
# LLM commander as its own unit (warehouse-bridge.service), so the ExecStart below pins
# sim:=false llm:=false to keep this a nav2-only / safety-only bring-up and avoid
# double-starting gz sim or the llm/nav2 bridge (#156 cross-lane -> #127).
# This unit MOVES the real robots, so it BindsTo= the Emergency Guardian:
# BindsTo + After guarantees Nav2 is never active without Layer 1 safety active
# and — unlike Requires= — also stops Nav2 when the guardian terminates
# UNEXPECTEDLY (crash / clean exit), not only on an explicit stop (systemd.unit(5);
# doc12:80-84, .claude/rules/safety.md).
Description=Warehouse Nav2 bring-up (per-bot Nav2 + twist_mux + VirtualScan)
Documentation=https://github.com/lll-kkk-ryuya/miniature-warehouse-robotics/blob/main/docs/setup/jetson-deploy.md
After=warehouse-safety.service
BindsTo=warehouse-safety.service
PartOf=warehouse.target

[Service]
Type=simple
EnvironmentFile=-/etc/warehouse/warehouse.env
User=warehouse
Group=warehouse
# RPLiDAR A1 over USB (doc02:180); dialout = serial-access group. use_sim_time:=false
# on real hardware; traffic_mode / map from the env file (see warehouse.env.example).
SupplementaryGroups=dialout
ExecStart=/opt/warehouse/deploy/jetson/bin/ros-exec.sh ros2 launch warehouse_bringup bringup.launch.py use_sim_time:=false sim:=false llm:=false traffic_mode:=${WAREHOUSE_TRAFFIC_MODE} map:=${WAREHOUSE_MAP}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=warehouse.target
